# srv.service
[Unit]
Description=Go web server

[Service]
Type=simple
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/quotes
ExecStart=/home/exedev/quotes/srv/srv
Restart=always
RestartSec=5
Environment=HOME=/home/exedev
Environment=USER=exedev
StandardOutput=journal
StandardError=journal

# Honeycomb observability (optional)
# Set HONEYCOMB_API_KEY to enable tracing
EnvironmentFile=-/home/exedev/quotes/.env

[Install]
WantedBy=multi-user.target
